# contrib/gentoo/net-dialup/accel-pptp/files/pptpd-init
#!/sbin/runscript

depend() {
	need net
}

start() {
	ebegin "Starting accel-pptpd"
	start-stop-daemon --start --quiet --exec /usr/sbin/accel-pptpd -- -p /var/run/pptpd.pid ${PPTPD_OPTS}
	eend $?
}

stop() {
	ebegin "Stopping accel-pptpd"
	start-stop-daemon --stop --quiet --pidfile /var/run/pptpd.pid
	result=$?
	start-stop-daemon --stop --quiet --oknodo --exec /usr/sbin/accel-pptpd
	result=$(( $result + $? ))
	eend $result
}
